# lens.service
[Unit]
Description=Gnome Lens Vector Search Daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/flatpak run org.gnome.Lens
Restart=on-failure
RestartSec=5
Nice=10
IOSchedulingClass=idle
StandardOutput=append:%h/.local/state/gnome-lens/daemon.log
StandardError=append:%h/.local/state/gnome-lens/daemon.log

[Install]
WantedBy=default.target